# autostart/dashboard.desktop
[Desktop Entry]
Type=Application
Name=Dashboard
Exec=/home/pi/bin/dashboard
