# autostart/telemetry.desktop
[Desktop Entry]
Type=Application
Name=Telemetry
Exec=telemetry
Comment[en_US]=Start Metis Telemetry System
