# autostart/telemetry.desktop
[Desktop Entry]
Type=Application
Name=Telemetry
Exec=/usr/local/bin/telemetry
