# autostart/xgpsspeed.desktop
[Desktop Entry]
Type=Application
Name=Speedometer
Exec=/usr/bin/xgpsspeed --speedunits=knots --maxspeed=20
